# docker/laravel.ini
[program:queues]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --tries=3 --timeout=60 --sleep=3 --queue=default,high
autostart=true
autorestart=true
numprocs=4
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/queues.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stopasgroup=true
killasgroup=true
startsecs=0
stopwaitsecs=10
exitcodes=0

[program:scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan schedule:work --verbose --no-interaction
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stopasgroup=true
killasgroup=true
startsecs=5
stopwaitsecs=10
exitcodes=0
