# docker/laravel.ini
[program:queues]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --tries=3 --timeout=60 --sleep=3 --queue=default,high
autostart=true
autorestart=unexpected
numprocs=4
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/queues.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stopasgroup=true
killasgroup=true
startsecs=0
stopwaitsecs=10
exitcodes=0,2

[program:scheduler]
command=php /var/www/artisan schedule:run --verbose --no-interaction >> /var/log/scheduler.log 2>&1
autostart=true
autorestart=unexpected
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stopasgroup=true
killasgroup=true
startsecs=0
stopwaitsecs=10
exitcodes=0,2
